# recipes-graphics/remote-kvm/remote-kvm/remote-kvm.service
[Unit]
Description=Remote KVM

[Service]
ExecStart={bindir}/remote-kvm

[Install]
RequiredBy={SYSTEMD_DEFAULT_TARGET}
